# kiosk/upstart/kiosk.conf
start on desktop-start
stop on desktop-end

console output

respawn

script
  exec > ~/.config/kiosk.log 2>&1
  exec /opt/sg_kiosk/kiosk.sh | logger -t kiosk
end script
